# camera.service
[Unit]
Description=Camera Utilty

[Service]
ExecStart=/bin/sh /home/pi/pi_camera/start_camera.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
